# assets/phi-bootstrap-update.service
[Unit]
Description=PhiNAS Bootstrap Update

[Service]
Type=simple
ExecStart=/phi/phi-bootstrap-update
