# assets/wisnuc-bootstrap-update.service
[Unit]
Description=Wisnuc Bootstrap Update

[Service]
Type=simple
ExecStart=/wisnuc/wisnuc-bootstrap-update
